# packages/uncloudd/output/DEBIAN/postinst
#!/bin/bash
set -e

if [ "$1" = "configure" ]; then
    adduser --system --home /var/lib/uncloud --shell /usr/sbin/nologin --group uncloud || true
    systemctl daemon-reload
    systemctl enable uncloudd.service
    systemctl start uncloudd.service || true

    echo ""
    echo "================================================================"
    echo " WARNING: This release of uncloud contains breaking changes."
    echo " Please refer to the uncloud release notes before proceeding:"
    echo " https://github.com/psviderski/uncloud/releases"
    echo "================================================================"
    echo ""
fi

# debian.griffo.io subscription notice (remove after the 1 Oct 2026 rollout settles)
if [ "$1" = "configure" ] && [ ! -e /etc/apt/auth.conf.d/debian.griffo.io.conf ]; then
    echo ""
    echo "debian.griffo.io: starting 1 October 2026 (noon UK), downloading packages"
    echo "from this repository requires a yearly subscription."
    echo "  How to subscribe: https://debian.griffo.io"
    echo "  Former Buy Me a Coffee supporters get a free lifetime subscription via Discord."
    echo ""
fi
